# openwrt-files/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

mkdir -p /mnt/usb
if [ -b /dev/sda1 ];
then
	mount /dev/sda1 /mnt/usb
	#make sure script is enables so it will stop process on shutdown/reboot
	/etc/init.d/logwifiids enable
	/etc/init.d/logwifiids start
fi

exit 0
